# ShellScript/myownservice.service
[Service]
ExecStart=/home/savera/Downloads/LinuxExercise/ShellScript/2_myownservice.sh
